# .config/systemd/user/podman-image-watch.path
[Unit]
Description=Watch for new Podman OCI image uploads

[Path]
PathModified=%h/.local/share/podman/incoming-images
PathExistsGlob=%h/.local/share/podman/incoming-images/*.tar

[Install]
WantedBy=default.target